# etc_salt/states/syncthing/supervisor.conf
[program:syncthing]
user = yann
environment = HOME="/home/yann"
command = syncthing
stopasgroup = true
killasgroup = true
